# eorchestra/roles/agora-tools/files/root_fs/etc/supervisor/conf.d/pynotify.conf
[program:install_certs]
command=/usr/local/bin/install_certs
stdout_events_enabled=true
stderr_events_enabled=true